# contrib/init/bitcoind.conf
description "Radiant Node Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env RADIANTD_BIN="/usr/bin/radiantd"
env RADIANTD_USER="radiant"
env RADIANTD_GROUP="radiant"
env RADIANTD_PIDDIR="/var/run/radiantd"
# upstart can't handle variables constructed with other variables
env RADIANTD_PIDFILE="/var/run/bitcoind/radiantd.pid"
env RADIANTD_CONFIGFILE="/etc/bitcoin/radiant.conf"
env RADIANTD_DATADIR="/var/lib/radiantd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # bitcoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$RADIANTD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run bitcoind."
        echo "The setting must appear in $RADIANTD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Bitcoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$RADIANTD_PIDDIR"
    chmod 0755 "$RADIANTD_PIDDIR"
    chown $RADIANTD_USER:$RADIANTD_GROUP "$RADIANTD_PIDDIR"
    chown $RADIANTD_USER:$RADIANTD_GROUP "$RADIANTD_CONFIGFILE"
    chmod 0660 "$RADIANTD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$RADIANTD_PIDFILE" \
    --chuid $RADIANTD_USER:$RADIANTD_GROUP \
    --exec "$RADIANTD_BIN" \
    -- \
    -pid="$RADIANTD_PIDFILE" \
    -conf="$RADIANTD_CONFIGFILE" \
    -datadir="$RADIANTD_DATADIR" \
    -disablewallet \
    -daemon

